# desktop/cron/crontab
0 12 * * * bash ~/dotfiles/cron/backup.sh >/tmp/crontab-stdout.log 2>/tmp/crontab-stderr.log2
